# extra/elering2mqtt.service
[Unit]
Description=elering2mqtt

[Service]
Type=oneshot
ExecStart=/usr/local/bin/elering2mqtt
