# extra/ow2mqtt.service
[Unit]
Description=ow2mqtt

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ow2mqtt
